# LPIC2-201/objetivos/inicio_sistema/Systemd_service/system-syncrorsync.timer
[Unit]
Description=sincroniza sincro con rsync cada minuto

[Timer]
OnUnitActiveSec=1m
Unit=system-sincrorsync.service

[Install]
WantedBy=timers.target
